# scripts/com.lotterypython.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lotterypython.daily</string>

    <key>Comment</key>
    <string>LotteryPython Daily Automation - Runs at 21:30 every day</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/Bruce/git/LotteryPython && source .venv/bin/activate && python scripts/daily_automation.py --type auto 2>&1 | tee -a logs/launchd.log</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/Bruce/git/LotteryPython</string>

    <key>StandardOutPath</key>
    <string>/Users/Bruce/git/LotteryPython/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Bruce/git/LotteryPython/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/Bruce/git/LotteryPython</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
